# usr/lib/systemd/system/qubes-updates-proxy-forwarder.socket.d/qubes-updates-cache.conf
[Unit]
ConditionPathExists=!/var/run/qubes-service/qubes-updates-cache
